# .config/autostart/quiterss.desktop
[Desktop Entry]
Type=Application
Exec=quiterss
Hidden=false
Name[en_US]=QuiteRSS
Name=QuiteRSS
Comment[en_US]=
Comment=
X-MATE-Autostart-Delay=0
